# systemd/authem-sync.timer
[Unit]
Description=Authem configuration syncer

[Timer]
OnBootSec=5min
OnCalendar=*:0/5
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
