# configs/systemd/authem-sync.service
[Unit]
Description=Authem configuration synchronization

[Service]
Type=oneshot
ExecStart=/usr/bin/authem-sync
